# etc/systemd/system/mpd-led.service
[Unit]
Description=mpd gpio led

[Service]
ExecStart=/home/alice/bin/mpd-led.py

[Install]
WantedBy=multi-user.target
