# etc/systemd/system/mpd-button.service
[Unit]
Description=mpd gpio button

[Service]
ExecStart=/home/alice/bin/mpd-button.py

[Install]
WantedBy=multi-user.target
